# com.user.alttabrestarter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.alttabrestarter</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/username/path/to/alttab-restarter.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/alttabrestarter.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/alttabrestarter.err</string>
</dict>
</plist>
